# setup/settings/profile.d/alias.sh
#!/bin/sh

alias phoebus=~/epics-tools/setup/startup/start_phoebus.sh
alias css=~/epics-tools/setup/startup/start_phoebus.sh
