# configure/profile.d/configure_jvm_env.sh
#!/bin/sh

export JAVA_HOME=/usr/lib/jvm/jdk-11.0.1
export PATH="$JAVA_HOME/bin:$PATH"

export MVN_HOME=/usr/lib/apache-maven-3.6.0
export PATH="$MVN_HOME/bin:$PATH"
