# configure/profile.d/configure_jvm_env.sh
#!/bin/sh

export JAVA_HOME=/opt/lib/jvm/jdk-11.0.2
export PATH="$JAVA_HOME/bin:$PATH"

export MVN_HOME=/opt/lib/apache-maven-3.6.0
export PATH="$MVN_HOME/bin:$PATH"
